# docker/sphinx/files/etc/supervisord.d/cron_log.conf
[program:cron_log]
command = tail --retry -f /var/log/cron.log
redirect_stderr=true
