# docker/sphinx/files/etc/supervisord.d/rsyslogd.conf
[program:rsyslogd]
command = rsyslogd -nd
redirect_stderr=true

