# docker/sphinx/files/etc/supervisord.d/sphinx.conf
[program:sphnix]
command = searchd --nodetach --config /usr/local/etc/sphinx.conf 
redirect_stderr = true
