# docker/sharetribe/files/etc/supervisord.d/configure.conf
[program:configure]
command = bundle exec rake ts:configure
redirect_stderr = true
directory = /opt/app
autostart = false
