# docker/sharetribe/files/etc/supervisord.d/copy_priv_2_pub.conf
[program:copy_priv_2_pub]
command =  /bin/bash -c "cp -R /opt/app/private/* /opt/app/public"
redirect_stderr = true
autostart = false
