# docker/sharetribe/files/etc/supervisord.d/cron_log.conf
[program:cron_log]
command = tail -f /var/log/cron.log
autostart = false
