# docker/sharetribe/files/etc/supervisord.d/dbinit.conf
[program:dbinit]
environment = DISABLE_DATABASE_ENVIRONMENT_CHECK=1
command =  bundle exec rake db:create db:structure:load
redirect_stderr = true
autostart = false

