# docker/sharetribe/files/etc/supervisord.d/index.conf
[program:index]
command = bundle exec rake ts:index
redirect_stderr = true
user = app
directory = /opt/app
autostart = false
