# docker/sharetribe/files/etc/supervisord.d/nginx.conf
[program:sphnix]
command = nginx -g 'daemon off;'
redirect_stderr = true

