# docker/sharetribe/files/etc/supervisord.d/sphinx.conf
[program:sphnix]
command = searchd --nodetach --config /sharetribe/data/sphinx/production/sphinx.conf 
redirect_stderr = true
user = app
