# docker/sharetribe/files/etc/supervisord.d/sphinx.conf
[program:sphnix]
command = searchd --nodetach --config /extvol/sphinx/production/sphinx.conf 
redirect_stderr = true

