# docker/sharetribe/files/etc/supervisord.d/web.conf
[program:web]
command = bundle exec rails server
redirect_stderr = true
directory = /opt/app
autostart = false
