# docker/sharetribe/files/etc/supervisord.d/worker.conf
[program:worker]
command = bundle exec rake jobs:work
redirect_stderr = true
user = app
directory = /opt/app
